# templates/agenticsec-log-cleanup.service.template
[Unit]
Description=AgenticSec Edge Log Cleanup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/agenticsec-log-cleanup.sh
